# watch-dog.timer
[Unit]
Description=Timer for check-services

[Timer]
OnBootSec=5s
OnUnitActiveSec=5s
Unit=watch-dog.service

[Install]
WantedBy=timers.target