# watch-dog.timer
[Unit]
Description=Timer for check-services

[Timer]
OnBootSec=15s
OnUnitActiveSec=15s
Unit=watch-dog.service

[Install]
WantedBy=timers.target